# systemd/trojita-google-contacts.service
[Unit]
Description=one-way sync from Google Contacts to Trojita

[Service]
Type=oneshot
ExecStart=@prefix/bin/trojita-google-contacts
